# deploy/hermes-paper.service
[Unit]
Description=Hermes v2 Paper — Docker Compose stack (bot + dashboard + nginx)
Documentation=https://github.com/minh99085/bot-3
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/financial-freedom-bot
EnvironmentFile=-/opt/financial-freedom-bot/.env

# Prefer docker compose plugin; fall back to docker-compose binary
ExecStart=/bin/bash -c 'if docker compose version >/dev/null 2>&1; then docker compose up -d --build --remove-orphans; else docker-compose up -d --build --remove-orphans; fi'
ExecStop=/bin/bash -c 'if docker compose version >/dev/null 2>&1; then docker compose down; else docker-compose down; fi'
ExecReload=/bin/bash -c 'if docker compose version >/dev/null 2>&1; then docker compose up -d --build --remove-orphans; else docker-compose up -d --build --remove-orphans; fi'

TimeoutStartSec=600
TimeoutStopSec=120

[Install]
WantedBy=multi-user.target
